# build/scripts/overlays/etc/systemd/system/zephyrus-gaming-qos.service
[Unit]
Description=Zephyrus Gaming QoS (Traffic Prioritization)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/zephyrus-gaming-qos
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
