# build/scripts/overlays/etc/systemd/system/zephyrus-crimson-setup.service
[Unit]
Description=Zephyrus Crimson OS Setup
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/zephyrus-crimson-setup
RemainAfterExit=yes

[Install]
WantedBy=graphical-session.target
