# build/scripts/overlays/etc/systemd/system/zephyrus-probe-hardware.service
[Unit]
Description=Zephyrus OS Hardware Probe (DPTF + EC)
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/zephyrus-probe-hardware
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
