# config/layered-fixes/configs/etc/udev/rules.d/99-zephyrus-ac-governor.rules
# Trigger CPU governor switch on AC adapter connect/disconnect
SUBSYSTEM=="power_supply", ATTR{name}=="ACAD", ATTR{online}=="0", RUN+="/usr/local/bin/zephyrus-ac-governor.sh"
SUBSYSTEM=="power_supply", ATTR{name}=="ACAD", ATTR{online}=="1", RUN+="/usr/local/bin/zephyrus-ac-governor.sh"
